# contrib/init/friccoind.conf
description "Friccoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env FRICCOIND_BIN="/usr/bin/friccoind"
env FRICCOIND_USER="friccoin"
env FRICCOIND_GROUP="friccoin"
env FRICCOIND_PIDDIR="/var/run/friccoind"
# upstart can't handle variables constructed with other variables
env FRICCOIND_PIDFILE="/var/run/friccoind/friccoind.pid"
env FRICCOIND_CONFIGFILE="/etc/friccoin/friccoin.conf"
env FRICCOIND_DATADIR="/var/lib/friccoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # friccoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$FRICCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run friccoind."
        echo "The setting must appear in $FRICCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Friccoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$FRICCOIND_PIDDIR"
    chmod 0755 "$FRICCOIND_PIDDIR"
    chown $FRICCOIND_USER:$FRICCOIND_GROUP "$FRICCOIND_PIDDIR"
    chown $FRICCOIND_USER:$FRICCOIND_GROUP "$FRICCOIND_CONFIGFILE"
    chmod 0660 "$FRICCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$FRICCOIND_PIDFILE" \
    --chuid $FRICCOIND_USER:$FRICCOIND_GROUP \
    --exec "$FRICCOIND_BIN" \
    -- \
    -pid="$FRICCOIND_PIDFILE" \
    -conf="$FRICCOIND_CONFIGFILE" \
    -datadir="$FRICCOIND_DATADIR" \
    -disablewallet \
    -daemon

